# debsrc/dc2-boot_1.0-0/DEBIAN/prerm
#!/bin/sh
if [ "$1" = "remove" ]; then
  service dc2-boot stop
  update-rc.d -f dc2-boot remove
fi

