# debian/dc2-boot.upstart
description "DC2-boot NodeJS System Service"

start on filesystem and net-device-up IFACE=eth1
stop on runlevel [06]

setuid dc2-boot
setgid adm

console log
respawn

exec /usr/bin/nodejs /var/www/dc2-boot/code/index.js
