# debsrc/dc2-boot_1.0-0/etc/init/dc2-boot.conf
description "DC2-boot NodeJS System Service"

start on filesystem
stop on runlevel [06]

console output
respawn

exec /usr/bin/nodejs /var/www/dc2-boot/code/index.js
